# serpentine-copy.service
[Unit]
Description=Copy RNO-G data to DESY

